# ft_services/srcs/nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=sh -c "nginx -g 'daemon off;' && kill -s SIGTERM $(cat supervisor.pid)"
# autostart=true
# autorestart=true
startsecs=0
redirect_stderr=true

